# sw/klipper/DEBIAN/prerm
#!/bin/bash
for s in \
  klipper-dc \
  getty@tty1
do
  systemctl disable "$s"
  systemctl stop    "$s"
done

# this script may never fail
exit 0
