# sw/klipper/DEBIAN/prerm
#!/bin/bash
systemctl stop x11kiosk
systemctl disable x11kiosk

# this script may never fail
exit 0
