# sw/klipper/DEBIAN/prerm
#!/bin/bash
systemctl stop    getty@tty1
systemctl disable getty@tty1

# this script may never fail
exit 0
